# config/etc/systemd/system/netif@modem.service.d/override.conf
[Unit]
Before=ppp@modem.service
Wants=ppp@modem.service
