# config/etc/systemd/system/netif@lan.service.d/override.conf
[Service]
ExecStartPost=-systemd-run --no-block systemctl try-restart dhcp6c.service
ExecStartPost=-systemd-run --no-block systemctl try-restart kea.service
